# coderunner.conf
description "CodeRunner Service daemon"
author      "Shashank Kumar - shashankkumar.me@gmail.com"
 
# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown
 
# Automatically Respawn:
respawn
respawn limit 99 5
 
script
    # Not sure why $HOME is needed, but in the guide it was written as such:
#    export HOME="/root"
 
    exec CodeRunner 2>&1
end script
 
post-start script
   # Optionally put a script here that will notifiy CodeRunner has (re)started
end script
